# sen2agri-persistence/dist/sen2agri-persistence.service
[Unit]
Description=Persistence Layer for Sen2Agri

[Service]
BusName=org.esa.sen2agri.persistenceManager
ExecStart=/usr/local/bin/sen2agri-persistence
