# sen2agri-services/dist/sen2agri-services.service
[Unit]
Description=Services for Sen2Agri
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/share/sen2agri/sen2agri-services/bin/start.sh
User=sen2agri-service
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
