# sen2agri-orchestrator/dist/sen2agri-orchestrator.service
[Unit]
Description=Orchestrator for Sen2Agri

[Service]
BusName=org.esa.sen2agri.orchestrator
ExecStart=/usr/local/bin/sen2agri-orchestrator
