# sen2agri-processors/DEM-WB/test/dist/sen2agri-demmaccs.timer
[Unit]
Description=Sentinel-2 DEM/MACCS runner

[Timer]
Unit=sen2agri-demmaccs.service
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
