# sen2agri-processors/DEM-WB/test/dist/sen2agri-demmaccs.service
[Unit]
Description=Runs MACCS on L1C products
After=postgresql-9.4.service
Requires=postgresql-9.4.service

[Service]
Type=simple
ExecStart=/usr/share/sen2agri/sen2agri-demmaccs/demmaccs_launcher.py
User=sen2agri-service
