# sen2agri-processors/DEM-WB/test/dist/sen2agri-demmaccs.service
[Unit]
Description=Runs MACCS on L1C products
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/share/sen2agri/sen2agri-demmaccs/l2a_launcher.py
User=sen2agri-service
