# sen2agri-processors/DEM-WB/test/dist/sen2agri-demmaccs.service
[Unit]
Description=Runs MACCS on L1C products

[Service]
Type=simple
#ExecStart=/usr/bin/demmaccs_launcher.py
ExecStart=/usr/share/sen2agri/sen2agri-demmacs/demmaccs_launcher.py
User=sen2agri-service
