# sen2agri-scheduler/dist/sen2agri-scheduler.service
[Unit]
Description=Scheduler for Sen2Agri

[Service]
ExecStart=/usr/bin/sen2agri-scheduler
User=sen2agri-service
