# sen2agri-http-listener/dist/sen2agri-http-listener.service
[Unit]
Description=Sen2Agri HTTP Listener
After=postgresql-9.4.service
Requires=postgresql-9.4.service

[Service]
Type=simple
ExecStart=/usr/bin/sen2agri-http-listener
User=sen2agri-service
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
