# sen2agri-http-listener/dist/sen2agri-http-listener.service
[Unit]
Description=Sen2Agri HTTP Listener

[Service]
Type=simple
ExecStart=/usr/bin/sen2agri-http-listener

[Install]
WantedBy=multi-user.target
