# sen2agri-downloaders/dist/sen2agri-sentinel-downloader.service
[Unit]
Description=Downloads Sentinel-2 products
After=postgresql-9.4.service
Requires=postgresql-9.4.service

[Service]
Type=simple
ExecStart=/usr/share/sen2agri/sen2agri-downloaders/downloader.py -r s2 -s /usr/share/sen2agri/sen2agri-downloaders/apihub.txt -l amazon
#ExecStart=/usr/share/sen2agri/sen2agri-downloaders/downloader.py -r s2 -s /usr/share/sen2agri/sen2agri-downloaders/apihub.txt -l local -i "/eodata/Sentinel-2/MSI/L1C"
User=sen2agri-service
