# sen2agri-downloaders/dist/sen2agri-landsat-downloader.timer
[Unit]
Description=Landsat 8 downloader timer

[Timer]
Unit=sen2agri-landsat-downloader.service
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
