# sen2agri-downloaders/dist/sen2agri-landsat-downloader.service
[Unit]
Description=Downloads Landsat 8 products

[Service]
Type=simple
#ExecStart=/usr/bin/landsat_download.py -u /etc/sen2agri/usgs.txt -z yes
ExecStart=/usr/share/sen2agri/sen2agri-downloaders/downloader.py -r l8 -s /usr/share/sen2agri/sen2agri-downloaders/usgs.txt
User=sen2agri-service
