# sen2agri-downloaders/dist/sen2agri-sentinel-downloader.timer
[Unit]
Description=Sentinel-2 downloader timer

[Timer]
Unit=sen2agri-sentinel-downloader.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
