# sen2agri-executor/dist/sen2agri-executor.service
[Unit]
Description=Executor for Sen2Agri

[Service]
BusName=org.esa.sen2agri.processorsExecutor
ExecStart=/usr/bin/sen2agri-executor
User=sen2agri-service
RestartSec=5
Restart=on-failure
Before=slurmd.service

[Install]
WantedBy=multi-user.target
