# sen2agri-executor/dist/sen2agri-executor.service
[Unit]
Description=Executor for Sen2Agri
Before=slurmd.service
After=docker.service
Requires=docker.service

[Service]
BusName=org.esa.sen2agri.processorsExecutor
ExecStart=/usr/bin/sen2agri-executor
User=sen2agri-service
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
